# docker/develop/supervisord/conf.d/cobblerd.conf
[program:cobblerd]
command=cobblerd -F
redirect_stderr=true
