# docker/develop/supervisord/conf.d/dhcpd.conf
[program:dhcpd]
command=/usr/sbin/dhcpd -4 -f -cf /etc/dhcpd.conf -lf /var/lib/dhcp/db/dhcpd.leases -user dhcpd -group nogroup pxe
autostart=false
autorestart=false
redirect_stderr=true
