# docker/develop/supervisord/conf.d/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn cobbler.services:application
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
