# docker/develop/supervisord/conf.d/mongodb.conf
[program:mongo]
command=/usr/bin/mongod --bind_ip_all
autorestart=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
