# tests/setup_files/supervisord/conf.d/httpd.conf
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
redirect_stderr=true
