# docker/debs/Debian_10/supervisord/conf.d/apache2.conf
[program:apache2]
command=apachectl -D "FOREGROUND" -k start
redirect_stderr=true
