# docker/rpms/Fedora_34/supervisord/conf.d/dhcpd.conf
[program:dhcpd]
command=/usr/sbin/dhcpd -f -cf /etc/dhcp/dhcpd.conf -user dhcpd -group dhcpd --no-pid
redirect_stderr=true
autostart=false
