# docker/rpms/Fedora_34/supervisord/conf.d/tftpd.conf
[program:tftpd]
command=/usr/sbin/in.tftpd -L -u root -s /srv/tftpboot
redirect_stderr=true
