# docker/Fedora33/supervisord/conf.d/httpd.conf
[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
redirect_stderr=true
