# docker/rpms/opensuse_leap/supervisord/conf.d/dhcpd.conf
[program:dhcpd]
command=/usr/lib/dhcp/dhcpd -4 start
redirect_stderr=true
