# docker/rpms/opensuse_leap/supervisord/conf.d/httpd.conf
[program:httpd]
command=/usr/sbin/start_apache2 -DFOREGROUND -k start
redirect_stderr=true
