# docker/rpms/opensuse_leap/supervisord/conf.d/slapd.conf
[program:slapd]
command=/usr/lib/openldap/start
stderr_logfile=/var/log/supervisord/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisord/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
autostart=false
