# docker/rpms/opensuse_leap/supervisord/conf.d/slapd.conf
[program:slapd]
command=/usr/lib/openldap/start
redirect_stderr=true
